# driver/linux/10-sgx.rules
SUBSYSTEM=="sgx",KERNEL=="sgx",RUN+="/bin/chmod 666 /dev/$name"
SUBSYSTEM=="sgx_prv",KERNEL=="sgx",GROUP="sgx_prv",RUN+="/bin/chmod 660 /dev/$name"

