# tools/SGXPlatformRegistration/package/installer/deb/sgx-ra-service/sgx-ra-service-1.0/debian/postinst
#!/bin/sh

# Exit on error
set -e

systemctl daemon-reload
systemctl enable mpa_registration_tool.service
systemctl start mpa_registration_tool.service
